# config/norns-maiden.socket
[Unit]
Description=maiden socket

[Socket]
ListenStream=80
NoDelay=true
